# config/supervisor/polestar.ini
[program:hblog]
directory=/data/HBlog
command=/data/HBlog/venv/bin/uwsgi --ini uwsgi.ini
autostart=true
autorestart=unexpected
exitcodes=0,2
user=ymmapp
startsecs=5
startretries=3
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/data/HBlog/logs/hblog-uwsgi.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
